# systemd/carbonio-oidc.service
[Unit]
Description=Carbonio OIDC Connector (Encedo)
After=network.target

[Service]
Type=simple
User=zextras
Group=zextras
WorkingDirectory=/opt/zextras/oidc
ExecStart=/usr/bin/python3 /opt/zextras/oidc/server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
